# net-misc/matterjs-server/files/matterjs-server.service
[Unit]
Description=OHF Matter.js Server
Documentation=https://github.com/matter-js/matterjs-server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
DynamicUser=yes
User=matter-server
Group=matter-server
StateDirectory=matterjs-server
WorkingDirectory=/var/lib/matterjs-server
ExecStart=/usr/bin/matter-server --storage-path /var/lib/matterjs-server
Restart=always
RestartSec=5s

# Security hardening
ProtectSystem=strict
ReadWritePaths=/var/lib/matterjs-server
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
